# scripts/com.bip-buddy.vision.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bip-buddy.vision</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/admin/projects/bip-buddy/scripts/start-vision.sh</string>
        <string>start</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/admin/projects/bip-buddy</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/admin/projects/bip-buddy/logs/vision.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/admin/projects/bip-buddy/logs/vision.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
